# myrvm-integration/systemd/rvm-remote-camera.service
[Unit]
Description=RVM Remote Camera Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=my
Group=my
WorkingDirectory=/home/my/test-cv-yolo11-sam2-camera/myrvm-integration
ExecStart=/usr/bin/python3 services/remote_camera_service.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/home/my/test-cv-yolo11-sam2-camera/myrvm-integration
Environment=DISPLAY=:0

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/my/test-cv-yolo11-sam2-camera/myrvm-integration/logs
ReadWritePaths=/home/my/test-cv-yolo11-sam2-camera/myrvm-integration/storages

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rvm-remote-camera

[Install]
WantedBy=multi-user.target
